# Assignments/scraper_service&package/scraper-deb/DEBIAN/postinst
#!/bin/bash
cd /scraper_project
python3 -m venv venv
./venv/bin/pip install -r requirements.txt

systemctl daemon-reload
systemctl enable scraper.service
systemctl start scraper.service
